# kde/profile.d/slimbookos-plasma-default.sh
XDG_CONFIG_DIRS="/usr/share/slimbook-plasma-settings-xdg/":"$XDG_CONFIG_DIRS"
